# com.seokmogu.notion-meeting-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.seokmogu.notion-meeting-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>-m</string>
    <string>notion_meeting_sync</string>
    <string>serve</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/seokmogu/project/notion-meeting-sync</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/notion-meeting-sync.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/notion-meeting-sync.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NMS_NOTION_TOKEN</key>
    <string>YOUR_TOKEN_HERE</string>
    <key>NMS_WEBHOOK_SECRET</key>
    <string>YOUR_VERIFICATION_TOKEN_HERE</string>
    <key>NMS_GIT_REPO_PATH</key>
    <string>/Users/seokmogu/project/agentic-services-docs</string>
    <key>NMS_NOTION_DATABASE_ID</key>
    <string>2ef7d8322b04834aac8a8158796b0a9a</string>
    <key>NO_SSL_VERIFY</key>
    <string>1</string>
  </dict>
</dict>
</plist>
